# modules/openclaw-market-immersion/systemd/openclaw-market-immersion-midday.timer
[Unit]
Description=Schedule OpenClaw Market Immersion midday collection

[Timer]
OnCalendar=Mon..Fri *-*-* 12:15:00
Persistent=true
Unit=openclaw-market-immersion-midday.service

[Install]
WantedBy=timers.target
